# gateway/gateway.service
[Unit]
Description=AI Mission Control Gateway
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/mission-control/gateway
ExecStart=/opt/mission-control/gateway/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8400
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
