# scripts/nexaas-client-dashboard.service
[Unit]
Description=Nexmatic Client Dashboard
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nexaas/client-dashboard
ExecStart=/usr/bin/node node_modules/.bin/next start -p 3001
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexaas-client-dashboard
EnvironmentFile=/opt/nexaas/client-dashboard/.env.local
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
