# scripts/nexaas-terminal.service
[Unit]
Description=Nexmatic Terminal WebSocket Server
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nexaas/dashboard
ExecStart=/usr/bin/npx tsx lib/terminal-server.ts
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexaas-terminal
EnvironmentFile=/opt/nexaas/dashboard/.env.local
Environment=TERMINAL_PORT=3002
Environment=NEXAAS_ROOT=/opt/nexaas

[Install]
WantedBy=multi-user.target
